# tareas/vmartinez/20161007-init-script-file/vms-daemon-test.service
[Unit]
Description=Example of systemd service file

[Service]
Type=simple
ExecStart=/usr/bin/daemon
User=root
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multiuser.agent
